# nzbget.conf
[program:nzbget]
user = jobski
startsecs = 0
autorestart = False
command = /home/jobski/start.sh
